# crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
TZ=Europe/Berlin
DATA_DIR=/app/data

# Morning briefing at 07:00 Frankfurt time
0 7 * * * root cd /app && python main.py --mode daily >> /proc/1/fd/1 2>&1

# Poll every 15 min during waking hours (07:05–22:50 Frankfurt time)
5,20,35,50 7-22 * * * root cd /app && python main.py --mode poll >> /proc/1/fd/1 2>&1
