# crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
TZ=Europe/Berlin
DATA_DIR=/app/data

# Morning briefing at 06:45 Frankfurt time
45 6 * * * root cd /app && python main.py --mode daily >> /proc/1/fd/1 2>&1

# Poll every 10 min during waking hours (07:00–22:50 Frankfurt time)
0,10,20,30,40,50 7-22 * * * root cd /app && python main.py --mode poll >> /proc/1/fd/1 2>&1
